# bluegreen.service
[Unit]
Description=Catalog Service
After=network.target

[Service]
ExecStart=/opt/installCatalog.sh

[Install]
WantedBy=multi-user.target
